# supervisor/startS.conf
[program:startS]
command=/home/ubuntu/RCProj/startS.sh
autostart=true
autorestart=true
stderr_logfile=/home/ubuntu/RCProj/server.err.log
stdout_logfile=/home/ubuntu/RCProj/server.out.log
